# systemd-1_237-r0/bsh-end-of-line-test.target
[Unit]
Description=BSH End-of-Line Test Target
Requires=basic.target
After=basic.target
Wants=getty.target
AllowIsolate=yes
